# service/recv/run
#!/bin/sh

PORT=$(cat port)

exec gst-launch-1.0 udpsrc address=:: port=$PORT ! 'application/x-rtp,payload=96,encoding-name=H264'! rtph264depay ! decodebin ! glimagesink
