# service/send/run
#!/bin/sh -ex

# TODO: Move in root context
# If use raspicam
# !su -c 'sh modprobe bcm2835_v4l2'

cd $(dirname $0)

#HOST=$(<host)
#PORT=$(<port)
HOST=localhost
PORT=5555

LOW=200000
BITRATE=$LOW

VIDEO_CAPS="video/x-h264,profile=high"
#VIDEO_CAPS="video/x-h264,profile=baseline"
VIDEO_SOURCE="videotestsrc"

# If use raspicam
# !su -c 'sh modprobe bcm2835_v4l2'
GST_FLAGS="-f -vvv"


GST_SRC='v4l2src ! "video/x-raw,width=640,height=480" !'
GST_ENCODE='omxh264enc target-bitrate='"$BITRATE"' control-rate=1 ! video/x-h264,profile=baseline !'
GST_RTP='h264parse ! rtph264pay pt=96 config-interval=1 ! '
GST_SEND='udpsink host='$HOST' port='$PORT


exec sh -c "gst-launch-1.0 $GST_FLAGS $GST_SRC $GST_ENCODE $GST_RTP $GST_SEND"


# PC VERSION
#exec gst-launch-1.0 v4l2src ! "video/x-raw,width=640,height=480" ! x264enc ! rtph264pay ! udpsink host=localhost port=5555
